# cronjober/crons/set-camera-time.cron
20 4 * * *  python3 /scripts/onvif_camera_set_time.py | tee -a /var/log/cron.log 2>&1
