# cronjober/crontab
*/5 * * * * * * echo "hello from Supercronic"
